# packages/trimarlib-sysfsgpio/trimarlib_sysfsgpio-0.2.1-py3-none-any.whl/sysfsgpio/resources/sysfsgpio.rules
ACTION!="remove", \
  KERNEL=="gpio*", SUBSYSTEM=="gpio", \
  RUN+="/bin/sh -c 'chown -R :staff /sys%p'", \
  RUN+="/bin/sh -c 'chmod -R g+rw /sys%p'"
