# packages/vlab-ipam-api/vlab_ipam_api-0.0.1-py3-none-any.whl/vlab_ipam_api/vlab-ipam.service
[Unit]
Description=RESTful API for vLab IPAM service
After=network.target

[Service]
ExecStart=/todo/uwsgi --ini app.ini
RuntimeDirectory=/usr/local/lib/python3.6/dist-packages/vlab_ipam_api
Restart=always

[Install]
WantedBy=multi-user.target
